# packages/net-misc/miredo/files/isatapd.init.d-0
#!/sbin/runscript
# Copyright 2009 Timothy Redaelli 
# Distributed under the terms of the GNU General Public License v2

depend() {
	need net
	use dns logger
}

start() {
	ebegin "Starting isatapd"
	start-stop-daemon --start --quiet --exec /usr/sbin/isatapd -- \
		${ISATAPD_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping isatapd"
	start-stop-daemon --stop --quiet --pidfile /var/run/isatapd.pid
	eend $?
}
