# packages/net-misc/miredo/files/miredo.init.d-0
#!/sbin/runscript
# Copyright 2009 Timothy Redaelli 
# Distributed under the terms of the GNU General Public License v2

depend() {
	need net
	use dns logger
}

start() {
	ebegin "Starting Miredo"
	start-stop-daemon --start --quiet --exec /usr/sbin/miredo -- \
		${MIREDO_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping Miredo"
	start-stop-daemon --stop --quiet --pidfile /var/run/miredo.pid
	eend $?
}
